# nginx-supervisor.ini
[supervisord]
user = root
nodaemon=true

[program:php-fpm]
command = /usr/sbin/php-fpm7 -F
user = root
autostart = true
#autorestart=unexpected
redirect_stderr=true
exitcodes=1

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true
#autorestart=unexpected
redirect_stderr=true
exitcodes=1
